# deploy/ec2/hydraflow.service
[Unit]
Description=@SERVICE_DESCRIPTION@
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@WORKING_DIRECTORY@
EnvironmentFile=-@ENV_FILE@
Environment=HYDRAFLOW_HOME_DIR=@HYDRAFLOW_HOME_DIR@
RuntimeDirectory=@RUNTIME_DIRECTORY@
ExecStart=@EXEC_START@
Restart=always
RestartSec=5
TimeoutStopSec=60
KillSignal=SIGINT
StandardOutput=append:@LOG_PATH@
StandardError=append:@LOG_PATH@
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
